# packaging/serve-engine.service
# systemd unit for the serve-engine daemon.
#
# Install:
#   sudo cp packaging/serve-engine.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now serve-engine
#
# This unit assumes:
#   - A `serve` user exists (sudo useradd -r -s /bin/false serve).
#   - /var/lib/serve is owned by `serve:serve` and is set as
#     SERVE_HOME (see config.py).
#   - The serve-engine package is installed in a venv at /opt/serve/venv
#     (e.g. `uv venv /opt/serve/venv && /opt/serve/venv/bin/pip install
#     -e /opt/serve/serve-engine`).
#
# Adjust the ExecStart path to match your install layout.

[Unit]
Description=serve-engine inference router daemon
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=exec
User=serve
Group=serve
WorkingDirectory=/var/lib/serve
Environment=SERVE_HOME=/var/lib/serve
ExecStart=/opt/serve/venv/bin/serve daemon start --foreground
Restart=on-failure
RestartSec=5

# Journald captures stdout/stderr — no log file in $HOME.
StandardOutput=journal
StandardError=journal

# Modest hardening. Loosen these if your deployments need additional
# capabilities (e.g. ReadWritePaths if model weights live elsewhere).
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/serve /var/run/docker.sock

[Install]
WantedBy=multi-user.target
